# edge/scale-daemon/scale-daemon.service
[Unit]
Description=OmniGate Scale Daemon — TCP weight reader and ingestor forwarder
Documentation=https://github.com/TruckGuard/omnigate
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=omnigate
Group=omnigate

ExecStart=/usr/local/bin/scale-daemon --config /etc/omnigate/scale-daemon.json

Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/omnigate

StandardOutput=journal
StandardError=journal
SyslogIdentifier=scale-daemon

[Install]
WantedBy=multi-user.target
